# base/var/spool/cron/crontabs/root
0 1 * * 0   ligero-fullbackup-S3.sh
0 1 * * 1-6 ligero-partialbackup-S3.sh